# scripts/launchd/ai.openclaw.mlb.update-results.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.openclaw.mlb.update-results</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/duste/.openclaw/workspace-zongzhihui/mlb_betting_system/.venv/bin/python3</string>
        <string>/Users/duste/.openclaw/workspace-zongzhihui/mlb_betting_system/scripts/update_results.py</string>
        <string>--tracker</string>
        <string>/Users/duste/.openclaw/workspace-zongzhihui/mlb_betting_system/data/performance_tracker.csv</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/duste/.openclaw/workspace-zongzhihui/mlb_betting_system</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/duste/.openclaw/workspace-zongzhihui/mlb_betting_system/logs/update_results.launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/duste/.openclaw/workspace-zongzhihui/mlb_betting_system/logs/update_results.launchd.err.log</string>
</dict>
</plist>
